# sv/httpbin/run
#!/bin/bash

#sv start socklog-unix || exit 1

source /etc/envvars

exec 2>&1
exec gunicorn --bind=0.0.0.0:5000 httpbin:app
